# pi/docker-update.service
[Unit]
Description=Weekly Docker Container Update
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/update-docker.sh
# StandardOutput and StandardError can be seen via journalctl --user -u docker-update.service

[Install]
WantedBy=default.target
